# home/mrtg/root.crontab
*/3 * * * * /home/scripts/GenHomePage.sh /home/scripts/MasterSystemList.txt
